# docker/s6-overlay/s6-rc.d/nats/run
#!/command/with-contenv sh
exec nats-server -js -m 8222
