# docker/s6-overlay/s6-rc.d/nats/run
#!/command/execlineb -P
nats-server -js -m 8222
