# docker/s6-overlay/s6-rc.d/ingestion/run
#!/command/with-contenv sh
set -eu

until wget --spider -q http://127.0.0.1:8222/healthz; do
	sleep 1
done

exec /usr/local/bin/ingestion
