# docker/s6-overlay/s6-rc.d/ingestion/run
#!/command/execlineb -P
/usr/local/bin/ingestion
