# docker/s6-overlay/s6-rc.d/ingestion/run
#!/command/with-contenv sh
exec /usr/local/bin/ingestion
