# docker/s6-overlay/s6-rc.d/processor/run
#!/command/execlineb -P
/usr/local/bin/processor
