# docker/s6-overlay/s6-rc.d/nginx/run
#!/command/execlineb -P
# Nginx service - serves Web UI and proxies API/WebSocket
nginx -g "daemon off;"
